# crontab
00 00 * * * /bin/sh /opt/centos-security/update.sh
